# scheduler/com.signaldesk.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.signaldesk.pipeline</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/emilygao/LocalDocuments/Projects/signaldesk/.venv/bin/python</string>
    <string>/Users/emilygao/LocalDocuments/Projects/signaldesk/pipeline/run_pipeline.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>/Users/emilygao/LocalDocuments/Projects/signaldesk</string>
    <key>PATH</key>
    <string>/Users/emilygao/LocalDocuments/Projects/signaldesk/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/emilygao/LocalDocuments/Projects/signaldesk</string>

  <key>StandardOutPath</key>
  <string>/Users/emilygao/LocalDocuments/Projects/signaldesk/logs/launchd_stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/emilygao/LocalDocuments/Projects/signaldesk/logs/launchd_stderr.log</string>

  <key>RunAtLoad</key>
  <false/>

</dict>
</plist>
